# luci-app-eqos/root/etc/uci-defaults/40_luci-eqos
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@eqos[-1]
	add ucitrack eqos
	set ucitrack.@eqos[-1].init=eqos
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
