# luci-theme-argon-new/root/etc/uci-defaults/30_luci-theme-argon-new
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Argon-New=/luci-static/argon-new
	set luci.main.mediaurlbase=/luci-static/argon-new
	commit luci
EOF
exit 0
