# luci-app-mosdns/root/etc/hotplug.d/iface/99-mosdns
#!/bin/sh

[ -f "/var/mosdns_sleep.lock" ] && exit 0

[ "$ACTION" = ifup ] && /etc/init.d/mosdns restart
