# luci-theme-opentopd/root/etc/uci-defaults/30_luci-theme-opentopd
#!/bin/sh
uci batch <<-EOF
	set luci.themes.opentopd=/luci-static/opentopd
	set luci.main.mediaurlbase=/luci-static/opentopd
	commit luci
EOF
exit 0
