# luci-theme-opentomcat/root/etc/uci-defaults/30_luci-theme-opentomcat
#!/bin/sh
uci batch <<-EOF
	set luci.themes.opentomcat=/luci-static/opentomcat
	set luci.main.mediaurlbase=/luci-static/opentomcat
	commit luci
EOF
exit 0
