# luci-app-store/root/etc/uci-defaults/luci-app-store
#!/bin/sh

/etc/init.d/istore boot

rm -rf /tmp/luci-indexcache /tmp/luci-modulecache
