# luci-app-control-speedlimit/root/etc/hotplug.d/iface/12-speedlimit
#!/bin/sh

[ "$ACTION" = "ifup" ] || exit 0
[ "$INTERFACE" = "br-lan" ] || exit 0
/etc/init.d/speedlimit start
