# luci-theme-opentomato/root/etc/uci-defaults/30_luci-theme-opentomato
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.OpenTomato=/luci-static/opentomato
		set luci.main.mediaurlbase=/luci-static/opentomato
		commit luci
	EOF
fi

exit 0
