# luci-app-advanced/root/etc/uci-defaults/40_luci-fb
#!/bin/sh

ver=$(date +%s)
[ -f /usr/lib/lua/luci/view/filebrowser.htm ] && sed -i "s/@ver/$ver/g" /usr/lib/lua/luci/view/filebrowser.htm

rm -f /tmp/luci-indexcache
exit 0
