# luci-theme-tomato/root/etc/uci-defaults/30_luci-theme-tomato
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.tomato=/luci-static/tomato
		set luci.main.mediaurlbase=/luci-static/tomato
		commit luci
	EOF
fi

exit 0
