# luci-theme-design/root/etc/uci-defaults/30_luci-theme-design
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.Design=/luci-static/design
		set luci.main.mediaurlbase=/luci-static/design
		commit luci
	EOF
fi

exit 0