# luci-app-partexp/luci-app-partexp/root/etc/uci-defaults/40_luci-partexp
#!/bin/sh
chmod +x /etc/init.d/partexp  >/dev/null 2>&1
[  `uci -q get partexp.global` ]  || uci set partexp.global=global
rm -rf /tmp/luci-modulecache /tmp/luci-indexcache*
exit 0
