# luci-app-partexp/luci-app-partexp/root/etc/uci-defaults/zzz_luci-app-partexp
#!/bin/sh

chmod +x /usr/bin/partexp /usr/libexec/rpcd/partexp
/etc/init.d/rpcd restart
rm -f /tmp/luci-indexcache
exit 0
