# luci-app-gpsysupgrade/root/etc/uci-defaults/luci-gpsysupgrade
rm -rf /tmp/luci-modulecache /tmp/luci-indexcache*
