# luci-themedog/root/etc/uci-defaults/30_luci-theme.themedog
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci get luci.themes.ThemeDog >/dev/null 2>&1 || \
	uci batch <<-EOF
		set luci.themes.ThemeDog=/luci-static/themedog
		set luci.main.mediaurlbase=/luci-static/themedog
		commit luci
	EOF
fi

exit 0
