# luci-app-bypass/root/etc/uci-defaults/luci-bypass
uci batch <<EOF
	delete ucitrack.@bypass[-1]
	add ucitrack bypass
	set ucitrack.@bypass[-1].init=bypass
	commit ucitrack
	delete firewall.bypass
	set firewall.bypass=include
	set firewall.bypass.type=script
	set firewall.bypass.path=/var/etc/bypass.include
	set firewall.bypass.reload=1
	commit firewall
	chmod +x /etc/init.d/bypass /usr/share/bypass/* >/dev/null 2>&1
	rm -rf /tmp/luci-modulecache/
	killall -HUP rpcd 2>/dev/null
EOF

if [[ "$(grep "option type 'vless'" /etc/config/bypass)" || "$(grep "option type 'vmess'" /etc/config/bypass)" ]]; then
sed -i "s/option uuid/option vmess_id/" /etc/config/bypass
for i in {1..100}
do
	uci -q set bypass.@servers[$i].type='v2ray'
	[ "$(uci -q get bypass.@servers[$i].type)" == 'vless' ] && uci -q set bypass.@servers[$i].v2ray_protocol='vless' || uci -q set bypass.@servers[$i].v2ray_protocol='vmess'
done
uci commit bypass
fi

exit 0
