# luci-app-guest-wifi/root/etc/uci-defaults/luci-i18n-guest-wifi-zh-cn
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	set luci.languages.zh_cn='chinese'
	commit luci
EOF

exit 0